# docker/install_scripts/rc.local
#!/bin/bash
bash /home/bit-server/scripts/update.sh
exit 0